# feedapp.service
[Unit]
Description=Regulatory Signal Aggregator (feedapp)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ken
WorkingDirectory=/home/ken/feedapp
EnvironmentFile=/home/ken/feedapp/.env
ExecStart=/home/ken/feedapp/.venv/bin/python main.py --continuous
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
